# cmus-syncthing.timer
[Unit]
Description=Triggers cmus-syncthing.service


[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
