# share/cmus-syncthing.service
[Unit]
Description= Exports cmus playlists to a folder

[Service]
Type=oneshot
ExecStart=cmus-syncthing

[Install]
WantedBy=multi-user.target
